# init/rpgd.conf
# rpgd - role-playing game daemon
#
# Starts rpgd for rpgclient connections

start on stopped rc RUNLEVEL=[2345]

console output
respawn
script
	exec /usr/local/bin/rpgd -d
end script

